# laravel/docker/supervisord/schooinqueue.ini
[program:laravelqueue]
command=/usr/local/bin/php /var/www/html/artisan queue:work database --quiet --memory=64 --tries=3
autostart=true
autorestart=true
user=www-data